# src/main/deb/control/postinst
#!/bin/sh

# Reload systemd manager configuration
systemctl daemon-reload

# Enable and start the ACME Server service
systemctl enable acmeserver
systemctl start acmeserver
